# docker-autostart.openrc
#!/sbin/openrc-run

command="python3"
command_args="/opt/docker/docker-autostart/discover.py /opt/docker/ /opt/docker/services/"

depend(){
	need docker
	after sshd
}

start(){
	$command $command_args
}

stop(){
	$command $command_args -a down
}
